# xwf/gateway/deploy/roles/fbstatsd/files/cmd2ods.service
[Service]
User=root
PermissionsStartOnly=true
ExecStart=/opt/xwf/cmd2ods/cmd2ods.rb
Restart=always
RestartSec=10s
[Install]
WantedBy=default.target
