# orc8r/cloud/deploy/roles/third_party/graphite/templates/debian_carbon_cache.j2
#!/bin/sh

### BEGIN INIT INFO
# Provides:        carbon-cache
# Required-Start:    $remote_fs $syslog
# Required-Stop:    $remote_fs $syslog
# Default-Start:    2 3 4 5
# Default-Stop:        
# Short-Description:    carbon-cache server
### END INIT INFO

set -e

test -x {{ graphite_install_path }}/bin/carbon-cache.py || exit 0

umask 022
ulimit -n {{ graphite_service_carbon_open_files_limit }}

. /lib/lsb/init-functions

if [ "x$2" != "x2nd" ]; then
	sudo -u {{ graphite_user }} $0 $1 2nd
	exit 0
fi

export PATH="${PATH:+$PATH:}/usr/sbin:/sbin"
export PYTHONPATH="{{ graphite_install_path }}/lib:{{ graphite_install_path }}/webapp"

case "$1" in
  start)
    log_daemon_msg "Starting carbon-cache daemon" "carbon-cache" || true
    if {{ graphite_install_path }}/bin/carbon-cache.py start; then
        log_end_msg 0 || true
    else
        log_end_msg 1 || true
    fi
    ;;
  stop)
    log_daemon_msg "Stopping carbon-cache daemon" "carbon-cache" || true
    if {{ graphite_install_path }}/bin/carbon-cache.py stop; then
        log_end_msg 0 || true
    else
        log_end_msg 1 || true
    fi
    ;;

  restart)
    log_daemon_msg "Restarting carbon-cache daemon" "carbon-cache" || true
    {{ graphite_install_path }}/bin/carbon-cache.py stop
    if {{ graphite_install_path }}/bin/carbon-cache.py start; then
        log_end_msg 0 || true
    else
        log_end_msg 1 || true
    fi
    ;;

  status)
    {{ graphite_install_path }}/bin/carbon-cache.py status && exit 0 || exit $?
    ;;

  *)
    log_action_msg "Usage: /etc/init.d/carbon-cache {start|stop|restart|status}" || true
    exit 1
esac

exit 0
