# lte/gateway/deploy/roles/magma/files/systemd/magma_liagentd.service
# Copyright 2020 The Magma Authors.

# This source code is licensed under the BSD-style license found in the
# LICENSE file in the root directory of this source tree.

# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
#
[Unit]
Description=Magma li agent service
PartOf=magma@pipelined.service
After=magma@pipelined.service

[Service]
Type=simple
EnvironmentFile=/etc/environment
ExecStart=/usr/local/bin/liagentd
ExecStopPost=/usr/bin/env python3 /usr/local/bin/service_util.py liagentd
MemoryAccounting=yes
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=liagentd
User=root
Restart=always
RestartSec=5
LimitCORE=infinity
MemoryLimit=300M

[Install]
WantedBy=multi-user.target